# finance/tryton-server/files/trytond-cron.in
#!/bin/sh

# PROVIDE: trytond-cron
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# trytond-cron knobs for /etc/rc.conf are as follows.
#
# trytond-cron_enable	Run trytond-cron at startup (default =NO)
# trytond-cron_args		Command line flags for trytond-cron (defaults
#						to trytond setting)
# trytond-cron_conf		Location of the trytond-cron config file(s)
#						(defaults to trytond setting)
#						Please note that configurations given later
#						on the line may, in part, or entirely,
#						override the ones specified earlier.
# trytond-cron_db		The database(s) for trytond-cron to operate
#						on (defaults to trytond setting)
# trytond-cron_logconf	Location of the logging config file (defaults
#						to trytond setting)
#
# Please note that trytond-cron_enable, trytond-cron_conf,
# trytond-cron_db, and trytond-cron_logconf must be set to trytond-
# acceptable values; hence null string values are ignored below
# intentionally.

. /etc/rc.subr

name=trytond-cron
rcvar=trytond-cron_enable

# Get defaults from tryond settings.
# NOTE: Must go before load_rc_config of this starup script.
. %%PREFIX%%/etc/trytond.rc.defaults

load_rc_config ${name}

: ${trytond-cron_enable:=NO}
: ${trytond-cron_args="${trytond_args}"}
: ${trytond-cron_conf:="${trytond_conf}"}
: ${trytond-cron_db:="${trytond_db}"}
: ${trytond-cron_logconf:="${trytond_logconf}"}

command=%%PREFIX%%/bin/${name}
command_interpreter=%%PYTHON_CMD%%
pidfile=/var/run/${name}.pid

# Please ensure to terminate command_args with an ampersand '&' sign
# at all times. Otherwise, you will see this script hang at boot time.
command_args="-c ${trytond-cron_conf} --logconf ${trytond-cron_logconf} --pidfile ${pidfile} &"

desc="Cron server for the Tryton Business Software suite."

run_rc_command "$1"
