# finance/tryton-server/files/trytond.in
#!/bin/sh

# PROVIDE: trytond
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# trytond knobs for /etc/rc.conf are as follows.
#
# trytond_enable    Run trytond at startup (default =NO)
# trytond_args      Command line flags for trytond (defaults to
#                   --)
# trytond_conf      Location of the trytond config file

. /etc/rc.subr

name=trytond
rcvar=trytond_enable

load_rc_config ${name}
: ${trytond_enable:=NO}
: ${trytond_args=}
: ${trytond_conf=}
: ${trytond_msg:=The ${name} script does nothing sensible yet.}

start_cmd="${name}_start"
stop_cmd=":"

desc="Server of the Tryton Business Software suite."

trytond_start()
{
	echo "${trytond_msg}"
}

run_rc_command "$1"
