# finance/tryton-server/files/trytond.in
#!/bin/sh

# PROVIDE: trytond
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# trytond knobs for /etc/rc.conf are as follows.
#
# trytond_enable	Run trytond at startup (default =NO)
# trytond_args		Command line flags for trytond (defaults to
#					--coroutine)
# trytond_conf		Location of the trytond config file(s) (defaults to
#					%%PREFIX%%/etc/trytond.conf \
#					%%PREFIX%%/etc/trytond.d/*)
#					Please note that configurations given later on the
#					line may, in part, or entirely, override the ones
#					specified earlier.
# trytond_db		The database(s) for trytond-cron to operate on
#					(defaults to %%PREFIX%%/var/db/trytond/db)
# trytond_logconf	Location of the logging config file (defaults to
#					%%PREFIX%%/etc/trytond.logconf)
#
# Please note that trytond_enable, trytond_conf, trytond_db, and
# trytond_logconf must be set to trytond-acceptable values; hence null
# string values are ignored below intentionally.

. /etc/rc.subr

name=trytond
rcvar=trytond_enable

# load_rc_config trytond is done via this include, please see there.
. %%PREFIX%%/etc/trytond.rc.defaults

command=%%PREFIX%%/bin/${name}
command_interpreter=%%PYTHON_CMD%%
pidfile=/var/run/${name}.pid

command_args="-c ${trytond_conf} --logconf ${trytond_logconf} --pidfile ${pidfile} &"

desc="Server of the Tryton Business Software suite."

run_rc_command "$1"
