# finance/tryton-server/files/trytond-worker.in
#!/bin/sh

# PROVIDE: trytond-worker
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# trytond-worker knobs for /etc/rc.conf are as follows.
#
# trytond-worker_enable		Run trytond-worker at startup (default
#							=NO)
# trytond-worker_args		Command line flags for trytond-worker
#							(defaults to trytond setting)
# trytond-worker_conf		Location of the trytond-worker config
#							file(s) (defaults to trytond setting)
#							Please note that configurations given
#							later on the line may, in part, or
#							entirely, override the ones specified
#							earlier.
# trytond-worker_db			The database(s) for trytond-worker to
#							operate on (defaults to trytond setting)
# trytond-worker_logconf	Location of the logging config file
#							(defaults to trytond setting)
#
# Please note that trytond-worker_enable, trytond-worker_conf,
# trytond-worker_db, and trytond-worker_logconf must be set to
# trytond-acceptable values; hence null string values are ignored
# below intentionally.

. /etc/rc.subr

name=trytond-worker
rcvar=trytond-worker_enable

# Get defaults from tryond settings.
# NOTE: Must go before load_rc_config of this starup script.
. %%PREFIX%%/etc/trytond.rc.defaults

load_rc_config ${name}

: ${trytond-worker_enable:=NO}
: ${trytond-worker_args="${trytond_args}"}
: ${trytond-worker_conf:="${trytond_conf}"}
: ${trytond-worker_db:="${trytond_db}"}
: ${trytond-worker_logconf:="${trytond_logconf}"}

command=%%PREFIX%%/bin/${name}
command_interpreter=%%PYTHON_CMD%%
pidfile=/var/run/${name}.pid

command_args="-c ${trytond-worker_conf} --logconf ${trytond-worker_logconf} --pidfile ${pidfile}"

desc="Workers for the Tryton Business Software suite."

run_rc_command "$1"
